# glpi/etc/cron.d/glpi
*/1 * * * * www-data /usr/bin/php /usr/share/glpi/front/cron.php
